# playbooks/roles/setup_tor/templates/tor.service.j2
[Unit]
Description=Tor
Documentation=man:tor
Documentation=https://www.torproject.org/
After=network-online.target docker.socket firewalld.service
Wants=network-online.target

[Service]
ExecStartPre=/usr/bin/tor --verify-config
ExecStart=/usr/bin/tor
ExecReload=/bin/kill -HUP ${MAINPID}
TimeoutSec={{ (bridge_distribution | length > 0) | ternary(5, 30) }}
User=tor
Group=tor
Restart=on-failure
LimitNOFILE=65536
KillSignal=SIGINT

# Hardening
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_NET_BIND_SERVICE CAP_DAC_READ_SEARCH CAP_KILL
NoNewPrivileges={{ (bridge_distribution | length > 0) | ternary('no', 'yes') }}
PrivateDevices=no
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=full
ReadOnlyDirectories=/
ReadWriteDirectories=-/var/lib/tor
ReadWriteDirectories=-/var/log/tor

[Install]
WantedBy=multi-user.target
